# deploy/systemd/vpn-backend.service
[Unit]
Description=VPN Telegram SaaS Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpnbot
Group=vpnbot
WorkingDirectory=/opt/vpn-bot/current/backend-go
EnvironmentFile=/opt/vpn-bot/current/backend-go/.env
ExecStart=/opt/vpn-bot/current/backend-go/bin/backend
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vpn-bot/shared /opt/vpn-bot/current/backend-go

[Install]
WantedBy=multi-user.target
